# downstream/rootfs-overlay/etc/systemd/system/jackrose-safe-mode.target
[Unit]
Description=Jackrose Safe Mode Target
AllowIsolate=yes
