# downstream/rootfs-overlay/etc/systemd/system/jackrose-desktop-failure.target
[Unit]
Description=Jackrose Desktop Failure Target
AllowIsolate=yes
